# celery-worker-fixed.service
[Unit]
Description=Celery Worker for Trading Assistant
After=network.target redis-server.service

[Service]
Type=simple
User=elcrypto
Group=elcrypto
WorkingDirectory=/home/elcrypto/trading_assistant
Environment="PATH=/home/elcrypto/trading_assistant/venv/bin"

# Запускаем через celery_worker.py который использует Application Factory
ExecStart=/home/elcrypto/trading_assistant/venv/bin/celery \
    -A celery_worker:celery worker \
    --loglevel=info \
    --concurrency=2

Restart=on-failure
RestartSec=10

LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target