# celery-worker-fixed.service
[Unit]
Description=Celery Worker for Trading Assistant
After=network.target redis-server.service

[Service]
Type=simple
User=elcrypto
Group=elcrypto
WorkingDirectory=/home/elcrypto/trading_assistant
Environment="PATH=/home/elcrypto/trading_assistant/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/elcrypto/trading_assistant"

# Используем простой worker для загрузки задач
ExecStart=/home/elcrypto/trading_assistant/venv/bin/celery -A celery_app worker --loglevel=info --concurrency=2

# Перезапуск при сбое
Restart=on-failure
RestartSec=10s

# Лимиты
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
